# deploy/systemd/caller-worker.service
[Unit]
Description=Alchemyst caller worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/alchemyst-ai/may-2026/devops/quickstart/workers/caller-worker
EnvironmentFile=/etc/alchemyst/caller-worker.env
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
