# deploy/systemd/iii-engine.service
[Unit]
Description=iii engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/aws-distributed-inference/quickstart
ExecStart=/usr/local/bin/iii --config /opt/aws-distributed-inference/quickstart/config.aws.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
